# crontab
*/1 * * * * swuser echo cd /home/swuser; ./update.sh >> /home/swuser/cron.log 2>&1

